# build/kinc/etc/systemd/system/kinc-init.service
[Unit]
Description=kinc Kubernetes Cluster Initialization
After=crio.service
Requires=crio.service
# Phase 2: Config resolved by kinc-init.sh fallback logic
# Supports both mounted config (priority) and baked-in config (fallback)

[Service]
Type=oneshot
ExecStart=/etc/kinc/scripts/kinc-init.sh
RemainAfterExit=yes
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
