# runtime/quadlet/kinc-control-plane.container
[Unit]
Description=kinc Kubernetes Control Plane Container (Quadlet Backend)
Documentation=https://github.com/T0MASD/kinc
After=network-online.target kinc-var-data-volume.service kinc-config-volume.service
Wants=network-online.target kinc-var-data-volume.service kinc-config-volume.service

[Container]
Image=localhost/kinc/node:v1.33.5
ContainerName=kinc-control-plane
HostName=kinc-control-plane

# Security and privileges (rootless configuration)
PodmanArgs=--cap-add=SYS_ADMIN --cap-add=SYS_RESOURCE --cap-add=NET_ADMIN --cap-add=SETPCAP --cap-add=NET_RAW --cap-add=SYS_PTRACE --cap-add=DAC_OVERRIDE --cap-add=CHOWN --cap-add=FOWNER --cap-add=FSETID --cap-add=KILL --cap-add=SETGID --cap-add=SETUID --cap-add=NET_BIND_SERVICE --cap-add=SYS_CHROOT --cap-add=SETFCAP --cap-add=DAC_READ_SEARCH --cap-add=AUDIT_WRITE --device /dev/fuse --sysctl net.ipv6.conf.all.disable_ipv6=0 --sysctl net.ipv6.conf.all.keep_addr_on_down=1 --sysctl net.netfilter.nf_conntrack_tcp_timeout_established=86400 --sysctl net.netfilter.nf_conntrack_tcp_timeout_close_wait=3600

# Tmpfs mounts (rootless configuration)
Tmpfs=/tmp:rw,rprivate,nosuid,nodev,tmpcopyup
Tmpfs=/run:rw,rprivate,nosuid,nodev,tmpcopyup
Tmpfs=/run/lock:rw,rprivate,nosuid,nodev,tmpcopyup

# Volume mounts (rootless configuration)
Volume=kinc-var-data:/var
Volume=kinc-config:/etc/kinc/config:ro
# Share host container storage to cache images across all clusters
Volume=%h/.local/share/containers/storage:/root/.local/share/containers/storage:rw

# Network configuration
# Port must be sequential (6443, 6444, 6445...) because subnet IDs are derived from last 2 digits
PublishPort=127.0.0.1:CLUSTER_PORT:6443/tcp

# Environment variables
Environment=container=podman
Environment=KUBECONFIG=/etc/kubernetes/admin.conf

# Container startup command - use default entrypoint
# Exec line removed - using default ENTRYPOINT from container image

[Service]
Restart=always
RestartSec=10
TimeoutStartSec=900
TimeoutStopSec=180

[Install]
WantedBy=multi-user.target
